# sentryvision.service
[Unit]
Description=SentryVision Home Security System
Requires=docker.service
After=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/barun/Documents/home-security-non-docker
ExecStart=/usr/bin/docker-compose -f docker-compose.yml up -d
ExecStop=/usr/bin/docker-compose -f docker-compose.yml down
TimeoutStartSec=0
User=barun
Group=barun

# Restart policy
Restart=on-failure
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/barun/Documents/home-security-non-docker

[Install]
WantedBy=multi-user.target